# installer/linux/systemd/openitms.service
[Unit]
Description=OpenITMS-SMB core (fork of Semaphore UI)
After=openitms-db.service
Requires=openitms-db.service

[Service]
Type=simple
WorkingDirectory=@PREFIX@
Environment=QUICKWIN_PLUGINS_DIR=@PREFIX@/plugins
Environment=OPENITMS_PREFIX=@PREFIX@
Environment=QUICKWIN_CONFIG=@PREFIX@/config/config.json
Environment=QUICKWIN_CERTS_DIR=@PREFIX@/certs
# Gitea addr/org/token (install.sh sinh; rỗng nếu không bundle Gitea)
EnvironmentFile=-@PREFIX@/config/gitea.env
ExecStart=@PREFIX@/bin/@BIN@ server --config @PREFIX@/config/config.json
Restart=on-failure
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
